# timerbot.service
[Unit]
Description=EVE Online Timer Discord Bot
After=network.target

[Service]
Type=simple
User=timerbot
Group=timerbot
WorkingDirectory=/opt/timerbot
ExecStart=/usr/bin/python3 -m bot.main
Restart=always
RestartSec=60

# Logging
StandardOutput=append:/opt/timerbot/logs/bot.log
StandardError=append:/opt/timerbot/logs/bot.log

# Security
NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=/opt/timerbot/data /opt/timerbot/logs
PrivateTmp=true

[Install]
WantedBy=multi-user.target 